# modules/90images/special/crontab
## Double comments will be removed from the generated file.
## The generator will expand some variables.  So far the list is: $HOME
##
## m h  dom mon dow   command

# Rotate background images
6,21,36,51 * * * * DISPLAY=:0 $HOME/bin/setbgimage &> /dev/null
